# building/deb/DEBIAN/postinst
#!/bin/sh

set -e

invoke-rc.d nginx reload
invoke-rc.d php5-fpm restart
